# pkg/woys-mic.service
[Unit]
Description=woys persistent virtual microphone (woys-mic)
Documentation=https://github.com/alirexha/woys
# We need the user PipeWire daemon up; sound.target is its anchor.
After=pipewire.socket pipewire-pulse.socket sound.target
Requires=pipewire-pulse.socket
PartOf=sound.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/env woys pw setup
# assert the virtual mic is actually
# healthy after setup (woys-mic present, system default sink not a woys
# null-sink, no ALSA leak links). A failed ExecStartPost makes the unit
# `failed` instead of a green `active (exited)` that hides broken audio.
ExecStartPost=/usr/bin/env woys chain status --check
ExecStop=/usr/bin/env woys pw teardown
# Service is idempotent — restart cleanly on failure.
Restart=on-failure
RestartSec=2

[Install]
# Enable per-user; this is a `systemctl --user` unit.
WantedBy=default.target
